# backend/gunicorn.service
[Unit]
Description=Gunicorn daemon for Django project
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/logistic-crm/backend
Environment="PATH=/opt/logistic-crm/backend/venv/bin"
ExecStart=/opt/logistic-crm/backend/venv/bin/gunicorn \
    --config /opt/logistic-crm/backend/gunicorn_config.py \
    logistic_crm.wsgi:application

[Install]
WantedBy=multi-user.target 